# docker/root/etc/crontabs/ksl
* 8-22 * * * python3 /app/ksl-notifications/src/KSL-Items/getItemList.py 2>&1