# docker/root/etc/crontabs/ksl_file
* 8-22 * * * python3 /app/ksl-notifications/src/KSL-Items/getItemList.py 2>&1 > /config/ksl.log
